# docker/etc/supervisor.d/async.ini
[program:async]
command=/usr/bin/sh -c "./manage.py async"
autorestart=true
directory=/var/local/mailur/src
environment=PATH="/var/local/mailur/env/bin"
