# docker/etc/supervisor.d/uwsgi.ini
[program:uwsgi]
command=/usr/bin/uwsgi --ini /etc/uwsgi/uwsgi.ini
autorestart=false
