# docker/etc/supervisor.d/watcher.ini
[program:watcher]
command=/usr/bin/sh -c "echo $PATH; while inotifywait -e modify -r .; do ./manage.py touch; done"
autorestart=true
directory=/var/local/mailur/src
environment=PATH="/var/local/mailur/env/bin:$PATH"
