# deploy/supervisor.ini
[program:postgres]
command=/usr/bin/postgres -D /var/lib/postgres/data
user=postgres
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/postgres.log

[program:web]
command=/bin/sh -c "gunicorn app:web -c deploy/web.conf.py"
directory=/var/local/mailur/src
environment=PYTHONPATH=".:./deploy",PATH="../env/bin/"
user=http
group=http
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/app.log

[program:async]
command=/bin/sh -c "gunicorn app:async -c deploy/async.conf.py"
directory=/var/local/mailur/src
environment=PYTHONPATH=".:./deploy",PATH="../env/bin/"
user=http
group=http
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/app.log

[program:watcher]
command=/usr/bin/sh -c "while inotifywait -e modify -r .; do ./m touch; done"
directory=/var/local/mailur/src
environment=PATH="/bin:./node_modules/.bin"
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/watcher.log
