# deploy/supervisor.ini
[program:postgres]
command=/usr/bin/postgres -D /var/lib/postgres/data
autorestart=false
user=postgres
redirect_stderr=true

[program:web]
command=/var/local/mailur/env/bin/gunicorn app:wsgi -c deploy/web.conf.py
autorestart=false
directory=/var/local/mailur/src
user=http
group=http
stdout_logfile=/var/log/app.log
redirect_stderr=true

[program:async]
command=/var/local/mailur/env/bin/gunicorn app:async -c deploy/async.conf.py
autorestart=false
directory=/var/local/mailur/src
user=http
group=http
stdout_logfile=/var/log/app.log
redirect_stderr=true

[program:watcher]
command=/usr/bin/sh -c "while inotifywait -e modify -r .; do ./m touch; done"
autorestart=true
directory=/var/local/mailur/src
redirect_stderr=true
