# deploy/etc/supervisor.d/async.ini
[program:async]
command=/var/local/mailur/env/bin/gunicorn app:async -c deploy/async.conf.py
autorestart=false
directory=/var/local/mailur/src
user=http
group=http
stdout_logfile=/var/log/app.log
redirect_stderr=true
