# deploy/etc/supervisor.d/nginx.ini
[program:nginx]
command=/usr/bin/nginx
autorestart=false
redirect_stderr=true
