# deploy/etc/supervisor.d/postgres.ini
[program:postgres]
command=/usr/bin/postgres -D /var/lib/postgres/data
autorestart=false
user=postgres
redirect_stderr=true
