# deploy/etc/supervisor.d/postgres.ini
[program:postgres]
command=/usr/bin/postgres -D /var/lib/postgres/data
user=postgres
autorestart=false
