# deploy/etc/supervisor.d/watcher.ini
[program:watcher]
command=/usr/bin/sh -c "while inotifywait -e modify -r .; do ./m touch; done"
autorestart=true
directory=/var/local/mailur/src
redirect_stderr=true
